# Webhook/webhook.service
[Unit]
Description=Webhooks

[Service]
ExecStart=webhook -hooks /root/Homelab/Webhook/hooks.json -hotreload

[Install]
WantedBy=multi-user.target