# Webhook/webhook.service
[Unit]
Description=Webhooks

[Service]
ExecStart=webhook -hooks /etc/Homelab/Webhook/hooks.json -template -hotreload -verbose

[Install]
WantedBy=multi-user.target