# Webhook/webhook.service
[Unit]
Description=Webhooks

[Service]
ExecStart=webhook -hooks /var/www/webhooks/hooks.json -hotreload

[Install]
WantedBy=multi-user.target