# Ansible/playbooks/templates/wireguard/wg-watchdog.service.j2
[Unit]
Description=WireGuard watchdog – {{ wireguard_interface | default('wg0') }}
After=network.target wg-quick@{{ wireguard_interface | default('wg0') }}.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wg-watchdog.sh