# Ansible/playbooks/templates/wireguard/wg-watchdog.timer.j2
[Unit]
Description=WireGuard watchdog timer – {{ wireguard_interface | default('wg0') }}

[Timer]
OnBootSec=60s
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target